# supervisor/command_center.conf
[program: command_center]
command = python3 command_center.py
directory = /var/command_center
autorestart = unexpected
