# supervisor/user-settings-service.conf
[program:user_settings]
command = python3 user_settings_service.py
directory = /var/command_center/services/user_settings
autorestart = unexpected
